# templates/certbot-deploy-hook.service.j2
[Unit]
Description=Run deploy hook after certificate renewal for {{ certbot_cert_name }}

[Service]
Type=oneshot
ExecStart={{ certbot_letsencrypt_config_dir }}/deploy-hook.sh
